# lib/systemd/system/onion-grater.service.d/30_cpfpy.conf
## Copyright (C) 2012 - 2018 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

[Service]
## Run /usr/lib/onion-grater-merger as root to avoid permission conflicts.
ExecStartPre=+/usr/lib/onion-grater-merger

ExecStart=
ExecStart=/usr/lib/onion-grater --debug --listen-interface eth1
