# lib/systemd/system/tor-controlport-filter.service.d/30_cpfpy.conf
## This file is part of Whonix.
## Copyright (C) 2012 - 2015 Patrick Schleizer <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

[Service]
EnvironmentFile=/usr/share/whonix-gw-network-conf/network_internal_ip.txt
ExecStart=
ExecStart=/usr/lib/tor-controlport-filter --debug --listen-address $network_internal_ip
